# conf/bastion/units/container-cleaner.timer
[Unit]
Description=container-cleaner-timer

[Timer]
OnCalendar=minutely

[X-Fleet]
MachineOf=container-cleaner.service
