# conf/bastion/units/deis-preseed.service
[Unit]
Description=deis-preseed

[Service]
EnvironmentFile=/etc/environment
ExecStart=-/run/deis/bin/preseed

[X-Fleet]
Global=true