# conf/bastion/units/deis-backup.timer
[Unit]
Description=deis-backup-timer

[Timer]
OnCalendar=hourly

[X-Fleet]
MachineOf=deis-backup.service
