# conf/bastion/units/image-cleaner.timer
[Unit]
Description=image-cleaner-timer

[Timer]
OnCalendar=daily

[X-Fleet]
MachineOf=image-cleaner.service
