# custom-units/app-restarter.timer
[Unit]
Description=deis-unit-restarter-timer

[Timer]
OnCalendar=daily

[X-Fleet]
MachineOf=%p.service
