# custom-units/blackhole.service
[Unit]
Description=block malicious ssh connections
After=sshd.service
Requires=sshd.service

[Service]
ExecStart=cat <<'Test' | bash \
          while true ; do \
            echo \
            "hello" \
            sleep 1 \
          done \
          Test

[X-Fleet]
MachineMetadata=bastion
