# custom-units/blackhole.service
# Unit file generated from generators/blackhole-service.sh
[Unit]
Description=block malicious ssh connections
After=sshd.service
Requires=sshd.service

[Service]
# Script encoded from scripts/blackhole.sh
ExecStart=/bin/sh -c 'echo "c2V0IC1lbyBwaXBlZmFpbAoKIyBDaGVjayB0byBzZWUgaWYgdGhlIHByb2dyYW0gaXMgcnVubmluZywgZXJyb3IgaWYgaXQgaXMKaWYgWyAtZiAvdmFyL3J1bi9ibGFja2hvbGUucGlkIF0gOyB0aGVuCiAgZWNobyAiQWNjb3JkaW5nIHRvIC92YXIvcnVuL2JsYWNraG9sZS5waWQsIHRoaXMgc2NyaXB0IGlzIGFscmVhZHkgcnVubmluZy4gRXhpdGluZy4iIDE+JjIgOwogIGV4aXQgMSA7CmZpIDsKCiMgVHJhY2sgdGhlIHByb2Nlc3MKZWNobyAkQkFTSFBJRCA+IC92YXIvcnVuL2JsYWNraG9sZS5waWQgOwoKIyBLZXlzCmV4cG9ydCBXSElURUxJU1RfS0VZPS9ibGFja2hvbGUvd2hpdGVsaXN0CmV4cG9ydCBCTEFDS0xJU1RfS0VZPS9ibGFja2hvbGUvd2hpdGVsaXN0CgojIEVuc3VyZSBuYW1lc3BhY2VzIGFyZSBjcmVhdGVkCmV0Y2RjdGwgbWtkaXIgJFdISVRFTElTVF9LRVkgfHwgdHJ1ZQpldGNkY3RsIG1rZGlyICRCTEFDS0xJU1RfS0VZIHx8IHRydWUKCmFkZF90b193aGl0ZWxpc3QoKXsKICBpcD0kMQogIGlmICEgaXNfd2hpdGVsaXN0ZWQgIiRpcCIgOyB0aGVuCiAgICBkYXRlIHwgZXRjZGN0bCBzZXQgIiRXSElURUxJU1RfS0VZIi8iJGlwIgogIGZpCn0KCmFkZF90b19ibGFja2xpc3QoKXsKICBpcD0kMQogIGlmICEgZXRjZGN0bCBnZXQgIiRCTEFDS0xJU1RfS0VZIi8iJGlwIiA+IC9kZXYvbnVsbCAyPiYxIDsgdGhlbgogICAgZGF0ZSB8IGV0Y2RjdGwgc2V0ICIkQkxBQ0tMSVNUX0tFWSIvIiRpcCIKICBmaQp9Cgppc193aGl0ZWxpc3RlZCgpewogIGlwPSQxCiAgaWYgISBldGNkY3RsIGdldCAiJFdISVRFTElTVF9LRVkiLyIkaXAiID4gL2Rldi9udWxsIDI+JjEgOyB0aGVuCiAgICBmYWxzZQogIGVsc2UKICAgIHRydWUKICBmaQp9CgojIFdoaXRlbGlzdCBmbGVldCBtYWNoaW5lcwooCiAgbGFzdCAtLXRpbWUtZm9ybWF0IG5vdGltZSAtYXdpIHwgZ3JlcCBwdHMgfCBjdXQgLWMzOC0KICBmbGVldGN0bCBsaXN0LW1hY2hpbmVzIC0tZmllbGRzPWlwIC1uby1sZWdlbmQKKSB8IHNvcnQgfCB1bmlxIHwgXAp3aGlsZSByZWFkIGlwIDsgZG8KICBhZGRfdG9fd2hpdGVsaXN0ICIkaXAiCmRvbmUKCiMgU2Nhbm5pbmcgUHJvZ3JhbQooCiAgL3Vzci9iaW4vam91cm5hbGN0bCAtZiAtbyBqc29uIHwgXAogIGdyZXAgLS1saW5lLWJ1ZmZlcmVkIHNzaGQgfCBcCiAgZ3JlcCAtLWxpbmUtYnVmZmVyZWQgLWUgIkRpc2Nvbm5lY3Rpbmc6IFRvbyBtYW55IGF1dGhlbnRpY2F0aW9uIGZhaWx1cmVzXHxGYWlsZWQgcGFzc3dvcmQgZm9yIiB8IFwKICBzZWQgLXUgLWUgInMvXi4qc3NoZEAuKjoyMi0vLyIgLWUgInMvOi4qJC8vIiB8IFwKICB3aGlsZSByZWFkIGlwIDsgZG8KICAgIGFkZF90b19ibGFja2xpc3QgIiRpcCIKICBkb25lICYgUkVBREVSX1BJRD0kIQogICgoIEpPVVJOQUxDVExfUElEPVJFQURFUl9QSUQrMSApKSA7CiAgX3Rlcm0oKSB7IGVjaG8gS2lsbGluZyBqb3VybmFsY3RsIHBpZCAkSk9VUk5BTENUTF9QSUQgMj4mMSA7CiAgICBraWxsICRKT1VSTkFMQ1RMX1BJRCA7CiAgICBleGl0IDAgOwogIH0gOwogIHRyYXAgX3Rlcm0gU0lHSU5UIFNJR1RFUk0gRVhJVCA7CiAgd2FpdCAkUkVBREVSX1BJRAopICYgU0NBTk5FUl9QSUQ9JCErMSA7CgoKIyBCbG9ja2VyIFByb2dyYW0KIyBXYXRjaGVzIGZvciBuZXcga2V5cyBpbiBldGNkIGFuZCBibG9ja3MgdGhlbQooCiAgZXRjZGN0bCBleGVjLXdhdGNoIC0tcmVjdXJzaXZlIC9ibGFja2hvbGUvYmxhY2tsaXN0IC0tIHNoIC1jICcKICAgIGlwPSQoZWNobyAkRVRDRF9XQVRDSF9LRVkgfCBzZWQgInMvLipcLy8vIikKICAgIGlmICEgaXNfd2hpdGVsaXN0ZWQgOyB0aGVuCiAgICAgIChpcCByb3V0ZSBhZGQgYmxhY2tob2xlICIkaXAiICYmIGVjaG8gImFkZGVkICRpcCB0byBibG9jayBsaXN0IikgfCB0cnVlCiAgICBmaQogICcgJiBXQVRDSEVSX1BJRD0kIQoKICBfdGVybSgpIHsgZWNobyBLaWxsaW5nIGV0ZGN0bCBwaWQgJFdBVENIRVJfUElEIDI+JjEKICAgIGtpbGwgJFdBVENIRVJfUElECiAgICBleGl0IDAKICB9CiAgdHJhcCBfdGVybSBTSUdJTlQgU0lHVEVSTSBFWElUCiAgd2FpdCAkV0FUQ0hFUl9QSUQKKSAmIEJMT0NLRVJfUElEPSQhKzEKCiMgRXhpc3RpbmcgS2V5IEJsb2NrZXIKIyBCbG9ja3MgYWxsIGV4aXN0aW5nIGtleXMgaW4gZXRjZApmb3Iga2V5IGluICQoZXRjZGN0bCBscyAvYmxhY2tob2xlL2JsYWNrbGlzdCkgOyBkbwogIGlwPSQoZWNobyAiJGtleSIgfCBzZWQgInMvLipcLy8vIikKICBpZiAhIGlzX3doaXRlbGlzdGVkIDsgdGhlbgogICAgKGlwIHJvdXRlIGFkZCBibGFja2hvbGUgIiRpcCIgJiYgZWNobyAiYWRkZWQgJGlwIHRvIGJsb2NrIGxpc3QiKSB8IHRydWUKICBmaQpkb25lCgojIFRlcm1pbmF0aW9uIFNjcmlwdApfdGVybSgpIHsKICBlY2hvICJLaWxsaW5nIEJsYWNraG9sZSBQcm9jZXNzZXMiICAyPiYxCiAga2lsbCAkU0NBTk5FUl9QSUQgPiAvZGV2L251bGwgMj4mMSB8fCB0cnVlCiAga2lsbCAkQkxPQ0tFUl9QSUQgPiAvZGV2L251bGwgMj4mMSB8fCB0cnVlCiAgcm0gLWYgL3Zhci9ydW4vYmxhY2tob2xlLnBpZAogIGV4aXQgMAp9CnRyYXAgX3Rlcm0gU0lHSU5UIFNJR1RFUk0gRVhJVAplY2hvICRCQVNIUElEID4gL3Zhci9ydW4vYmxhY2tob2xlLnBpZAp3YWl0ICRTQ0FOTkVSX1BJRAp3YWl0ICRCTE9DS0VSX1BJRAo=" | base64 -d | bash'

[X-Fleet]
MachineMetadata=bastion

