# custom-units/blackhole.service
# Unit file generated from generators/blackhole-service.sh
[Unit]
Description=block malicious ssh connections
After=sshd.service
Requires=sshd.service

[Service]
# Script encoded from scripts/blackhole.sh
ExecStart=/bin/sh -c 'echo "c2V0IC1lbyBwaXBlZmFpbAoKIyBDaGVjayB0byBzZWUgaWYgdGhlIHByb2dyYW0gaXMgcnVubmluZywgZXJyb3IgaWYgaXQgaXMKaWYgWyAtZiAvdmFyL3J1bi9ibGFja2hvbGUucGlkIF0gOyB0aGVuCiAgZWNobyAiQWNjb3JkaW5nIHRvIC92YXIvcnVuL2JsYWNraG9sZS5waWQsIHRoaXMgc2NyaXB0IGlzIGFscmVhZHkgcnVubmluZy4gRXhpdGluZy4iIDE+JjIgOwogIGV4aXQgMSA7CmZpIDsKCiMgVHJhY2sgdGhlIHByb2Nlc3MKZWNobyAkQkFTSFBJRCA+IC92YXIvcnVuL2JsYWNraG9sZS5waWQgOwoKIyBLZXlzCmV4cG9ydCBXSElURUxJU1RfS0VZPS9ibGFja2hvbGUvd2hpdGVsaXN0CmV4cG9ydCBCTEFDS0xJU1RfS0VZPS9ibGFja2hvbGUvd2hpdGVsaXN0CgojIEVuc3VyZSBuYW1lc3BhY2VzIGFyZSBjcmVhdGVkCmV0Y2RjdGwgbWtkaXIgJFdISVRFTElTVF9LRVkKZXRjZGN0bCBta2RpciAkQkxBQ0tMSVNUX0tFWQoKYWRkX3RvX3doaXRlbGlzdCgpewogIGlwPSQxCiAgaWYgISBpc193aGl0ZWxpc3RlZCAiJGlwIiA7IHRoZW4KICAgIGRhdGUgfCBldGNkY3RsIHNldCAiJFdISVRFTElTVF9LRVkiLyIkaXAiCiAgZmkKfQoKYWRkX3RvX2JsYWNrbGlzdCgpewogIGlwPSQxCiAgaWYgISBldGNkY3RsIGdldCAiJEJMQUNLTElTVF9LRVkiLyIkaXAiID4gL2Rldi9udWxsIDI+JjEgOyB0aGVuCiAgICBkYXRlIHwgZXRjZGN0bCBzZXQgIiRCTEFDS0xJU1RfS0VZIi8iJGlwIgogIGZpCn0KCmlzX3doaXRlbGlzdGVkKCl7CiAgaXA9JDEKICBpZiAhIGV0Y2RjdGwgZ2V0ICIkV0hJVEVMSVNUX0tFWSIvIiRpcCIgPiAvZGV2L251bGwgMj4mMSA7IHRoZW4KICAgIGZhbHNlCiAgZWxzZQogICAgdHJ1ZQogIGZpCn0KCiMgV2hpdGVsaXN0IGZsZWV0IG1hY2hpbmVzCigKICBsYXN0IC0tdGltZS1mb3JtYXQgbm90aW1lIC1hd2kgfCBncmVwIHB0cyB8IGN1dCAtYzM4LQogIGZsZWV0Y3RsIGxpc3QtbWFjaGluZXMgLS1maWVsZHM9aXAgLW5vLWxlZ2VuZAopIHwgc29ydCB8IHVuaXEgfCBcCndoaWxlIHJlYWQgaXAgOyBkbwogIGFkZF90b193aGl0ZWxpc3QgIiRpcCIKZG9uZQoKIyBTY2FubmluZyBQcm9ncmFtCigKICAvdXNyL2Jpbi9qb3VybmFsY3RsIC1mIC1vIGpzb24gfCBcCiAgZ3JlcCAtLWxpbmUtYnVmZmVyZWQgc3NoZCB8IFwKICBncmVwIC0tbGluZS1idWZmZXJlZCAtZSAiRGlzY29ubmVjdGluZzogVG9vIG1hbnkgYXV0aGVudGljYXRpb24gZmFpbHVyZXNcfEZhaWxlZCBwYXNzd29yZCBmb3IiIHwgXAogIHNlZCAtdSAtZSAicy9eLipzc2hkQC4qOjIyLS8vIiAtZSAicy86LiokLy8iIHwgXAogIHdoaWxlIHJlYWQgaXAgOyBkbwogICAgYWRkX3RvX2JsYWNrbGlzdCAiJGlwIgogIGRvbmUgJiBSRUFERVJfUElEPSQhCiAgKCggSk9VUk5BTENUTF9QSUQ9UkVBREVSX1BJRCsxICkpIDsKICBfdGVybSgpIHsgZWNobyBLaWxsaW5nIGpvdXJuYWxjdGwgcGlkICRKT1VSTkFMQ1RMX1BJRCAyPiYxIDsKICAgIGtpbGwgJEpPVVJOQUxDVExfUElEIDsKICAgIGV4aXQgMCA7CiAgfSA7CiAgdHJhcCBfdGVybSBTSUdJTlQgU0lHVEVSTSBFWElUIDsKICB3YWl0ICRSRUFERVJfUElECikgJiBTQ0FOTkVSX1BJRD0kISsxIDsKCgojIEJsb2NrZXIgUHJvZ3JhbQojIFdhdGNoZXMgZm9yIG5ldyBrZXlzIGluIGV0Y2QgYW5kIGJsb2NrcyB0aGVtCigKICBldGNkY3RsIGV4ZWMtd2F0Y2ggLS1yZWN1cnNpdmUgL2JsYWNraG9sZS9ibGFja2xpc3QgLS0gc2ggLWMgJwogICAgaXA9JChlY2hvICRFVENEX1dBVENIX0tFWSB8IHNlZCAicy8uKlwvLy8iKQogICAgaWYgISBpc193aGl0ZWxpc3RlZCA7IHRoZW4KICAgICAgKGlwIHJvdXRlIGFkZCBibGFja2hvbGUgIiRpcCIgJiYgZWNobyAiYWRkZWQgJGlwIHRvIGJsb2NrIGxpc3QiKSB8IHRydWUKICAgIGZpCiAgJyAmIFdBVENIRVJfUElEPSQhCgogIF90ZXJtKCkgeyBlY2hvIEtpbGxpbmcgZXRkY3RsIHBpZCAkV0FUQ0hFUl9QSUQgMj4mMQogICAga2lsbCAkV0FUQ0hFUl9QSUQKICAgIGV4aXQgMAogIH0KICB0cmFwIF90ZXJtIFNJR0lOVCBTSUdURVJNIEVYSVQKICB3YWl0ICRXQVRDSEVSX1BJRAopICYgQkxPQ0tFUl9QSUQ9JCErMQoKIyBFeGlzdGluZyBLZXkgQmxvY2tlcgojIEJsb2NrcyBhbGwgZXhpc3Rpbmcga2V5cyBpbiBldGNkCmZvciBrZXkgaW4gJChldGNkY3RsIGxzIC9ibGFja2hvbGUvYmxhY2tsaXN0KSA7IGRvCiAgaXA9JChlY2hvICIka2V5IiB8IHNlZCAicy8uKlwvLy8iKQogIGlmICEgaXNfd2hpdGVsaXN0ZWQgOyB0aGVuCiAgICAoaXAgcm91dGUgYWRkIGJsYWNraG9sZSAiJGlwIiAmJiBlY2hvICJhZGRlZCAkaXAgdG8gYmxvY2sgbGlzdCIpIHwgdHJ1ZQogIGZpCmRvbmUKCiMgVGVybWluYXRpb24gU2NyaXB0Cl90ZXJtKCkgewogIGVjaG8gIktpbGxpbmcgQmxhY2tob2xlIFByb2Nlc3NlcyIgIDI+JjEKICBraWxsICRTQ0FOTkVSX1BJRCA+IC9kZXYvbnVsbCAyPiYxIHx8IHRydWUKICBraWxsICRCTE9DS0VSX1BJRCA+IC9kZXYvbnVsbCAyPiYxIHx8IHRydWUKICBybSAtZiAvdmFyL3J1bi9ibGFja2hvbGUucGlkCiAgZXhpdCAwCn0KdHJhcCBfdGVybSBTSUdJTlQgU0lHVEVSTSBFWElUCmVjaG8gJEJBU0hQSUQgPiAvdmFyL3J1bi9ibGFja2hvbGUucGlkCndhaXQgJFNDQU5ORVJfUElECndhaXQgJEJMT0NLRVJfUElECg==" | base64 -d | bash'

[X-Fleet]
MachineMetadata=bastion

