# custom-units/blackhole.service
# Unit file generated from generators/blackhole-service.sh
[Unit]
Description=block malicious ssh connections
After=sshd.service
Requires=sshd.service

[Service]
# Script encoded from scripts/blackhole.sh
ExecStart=/bin/sh -c 'echo "IyEvdXNyL2xvY2FsL2VudiBiYXNoCnNldCAtZW8gcGlwZWZhaWwKCiMgQ2hlY2sgdG8gc2VlIGlmIHRoZSBwcm9ncmFtIGlzIHJ1bm5pbmcsIGVycm9yIGlmIGl0IGlzCmlmIFsgLWYgL3Zhci9ydW4vYmxhY2tob2xlLnBpZCBdIDsgdGhlbgogIGVjaG8gIkFjY29yZGluZyB0byAvdmFyL3J1bi9ibGFja2hvbGUucGlkLCB0aGlzIHNjcmlwdCBpcyBhbHJlYWR5IHJ1bm5pbmcuIEV4aXRpbmcuIiAxPiYyIDsKICBleGl0IDEgOwpmaSA7CgojIFRyYWNrIHRoZSBwcm9jZXNzCmVjaG8gJEJBU0hQSUQgPiAvdmFyL3J1bi9ibGFja2hvbGUucGlkIDsKCiMgS2V5cwpleHBvcnQgV0hJVEVMSVNUX0tFWT0vYmxhY2tob2xlL3doaXRlbGlzdApleHBvcnQgQkxBQ0tMSVNUX0tFWT0vYmxhY2tob2xlL3doaXRlbGlzdAoKIyBFbnN1cmUgbmFtZXNwYWNlcyBhcmUgY3JlYXRlZApldGNkY3RsIG1rZGlyICRXSElURUxJU1RfS0VZIHx8IHRydWUKZXRjZGN0bCBta2RpciAkQkxBQ0tMSVNUX0tFWSB8fCB0cnVlCgphZGRfdG9fd2hpdGVsaXN0KCl7CiAgaXA9JDEKICBpZiAhIGlzX3doaXRlbGlzdGVkICIkaXAiIDsgdGhlbgogICAgZGF0ZSB8IGV0Y2RjdGwgc2V0ICIkV0hJVEVMSVNUX0tFWSIvIiRpcCIKICBmaQp9CgphZGRfdG9fYmxhY2tsaXN0KCl7CiAgaXA9JDEKICBpZiAhIGV0Y2RjdGwgZ2V0ICIkQkxBQ0tMSVNUX0tFWSIvIiRpcCIgPiAvZGV2L251bGwgMj4mMSA7IHRoZW4KICAgIGRhdGUgfCBldGNkY3RsIHNldCAiJEJMQUNLTElTVF9LRVkiLyIkaXAiCiAgZmkKfQoKaXNfd2hpdGVsaXN0ZWQoKXsKICBpcD0kMQogIGlmICEgZXRjZGN0bCBnZXQgIiRXSElURUxJU1RfS0VZIi8iJGlwIiA+IC9kZXYvbnVsbCAyPiYxIDsgdGhlbgogICAgZmFsc2UKICBlbHNlCiAgICB0cnVlCiAgZmkKfQoKIyBXaGl0ZWxpc3QgZmxlZXQgbWFjaGluZXMKKAogIGxhc3QgLS10aW1lLWZvcm1hdCBub3RpbWUgLWF3aSB8IGdyZXAgcHRzIHwgY3V0IC1jMzgtCiAgZmxlZXRjdGwgbGlzdC1tYWNoaW5lcyAtLWZpZWxkcz1pcCAtbm8tbGVnZW5kCikgfCBzb3J0IHwgdW5pcSB8IFwKd2hpbGUgcmVhZCBpcCA7IGRvCiAgYWRkX3RvX3doaXRlbGlzdCAiJGlwIgpkb25lCgojIFNjYW5uaW5nIFByb2dyYW0KKAogIC91c3IvYmluL2pvdXJuYWxjdGwgLWYgLW8ganNvbiB8IFwKICBncmVwIC0tbGluZS1idWZmZXJlZCBzc2hkIHwgXAogIGdyZXAgLS1saW5lLWJ1ZmZlcmVkIC1lICJEaXNjb25uZWN0aW5nOiBUb28gbWFueSBhdXRoZW50aWNhdGlvbiBmYWlsdXJlc1x8RmFpbGVkIHBhc3N3b3JkIGZvciIgfCBcCiAgc2VkIC11IC1lICJzL14uKnNzaGRALio6MjItLy8iIC1lICJzLzouKiQvLyIgfCBcCiAgd2hpbGUgcmVhZCBpcCA7IGRvCiAgICBhZGRfdG9fYmxhY2tsaXN0ICIkaXAiCiAgZG9uZSAmIFJFQURFUl9QSUQ9JCEKICAoKCBKT1VSTkFMQ1RMX1BJRD1SRUFERVJfUElEKzEgKSkgOwogIF90ZXJtKCkgeyBlY2hvIEtpbGxpbmcgam91cm5hbGN0bCBwaWQgJEpPVVJOQUxDVExfUElEIDI+JjEgOwogICAga2lsbCAkSk9VUk5BTENUTF9QSUQgOwogICAgZXhpdCAwIDsKICB9IDsKICB0cmFwIF90ZXJtIFNJR0lOVCBTSUdURVJNIEVYSVQgOwogIHdhaXQgJFJFQURFUl9QSUQKKSAmIFNDQU5ORVJfUElEPSQhKzEgOwoKCiMgQmxvY2tlciBQcm9ncmFtCiMgV2F0Y2hlcyBmb3IgbmV3IGtleXMgaW4gZXRjZCBhbmQgYmxvY2tzIHRoZW0KKAogIHJlYWQgLXIgLWQgJycgU0NSSVBUIDw8LSdFT0YnCmlwPSQoZWNobyAiJEVUQ0RfV0FUQ0hfS0VZIiB8IHNlZCAicy8uKlwvLy8iKQppZiAhIGlzX3doaXRlbGlzdGVkIDsgdGhlbgogIChpcCByb3V0ZSBhZGQgYmxhY2tob2xlICIkaXAiICYmIGVjaG8gImFkZGVkICRpcCB0byBibG9jayBsaXN0IikgfCB0cnVlCmZpCkVPRgogIGV0Y2RjdGwgZXhlYy13YXRjaCAtLXJlY3Vyc2l2ZSAvYmxhY2tob2xlL2JsYWNrbGlzdCAtLSBzaCAtYyAiJFNDUklQVCIgJiBXQVRDSEVSX1BJRD0kIQoKICBfdGVybSgpIHsgZWNobyBLaWxsaW5nIGV0ZGN0bCBwaWQgJFdBVENIRVJfUElEIDI+JjEKICAgIGtpbGwgJFdBVENIRVJfUElECiAgICBleGl0IDAKICB9CiAgdHJhcCBfdGVybSBTSUdJTlQgU0lHVEVSTSBFWElUCiAgd2FpdCAkV0FUQ0hFUl9QSUQKKSAmIEJMT0NLRVJfUElEPSQhKzEKCiMgRXhpc3RpbmcgS2V5IEJsb2NrZXIKIyBCbG9ja3MgYWxsIGV4aXN0aW5nIGtleXMgaW4gZXRjZApmb3Iga2V5IGluICQoZXRjZGN0bCBscyAvYmxhY2tob2xlL2JsYWNrbGlzdCkgOyBkbwogIGlwPSQoZWNobyAiJGtleSIgfCBzZWQgInMvLipcLy8vIikKICBpZiAhIGlzX3doaXRlbGlzdGVkIDsgdGhlbgogICAgKGlwIHJvdXRlIGFkZCBibGFja2hvbGUgIiRpcCIgJiYgZWNobyAiYWRkZWQgJGlwIHRvIGJsb2NrIGxpc3QiKSB8IHRydWUKICBmaQpkb25lCgojIFRlcm1pbmF0aW9uIFNjcmlwdApfdGVybSgpIHsKICBlY2hvICJLaWxsaW5nIEJsYWNraG9sZSBQcm9jZXNzZXMiICAyPiYxCiAga2lsbCAkU0NBTk5FUl9QSUQgPiAvZGV2L251bGwgMj4mMSB8fCB0cnVlCiAga2lsbCAkQkxPQ0tFUl9QSUQgPiAvZGV2L251bGwgMj4mMSB8fCB0cnVlCiAgcm0gLWYgL3Zhci9ydW4vYmxhY2tob2xlLnBpZAogIGV4aXQgMAp9CnRyYXAgX3Rlcm0gU0lHSU5UIFNJR1RFUk0gRVhJVAplY2hvICRCQVNIUElEID4gL3Zhci9ydW4vYmxhY2tob2xlLnBpZAp3YWl0ICRTQ0FOTkVSX1BJRAp3YWl0ICRCTE9DS0VSX1BJRAo=" | base64 -d | bash'

[X-Fleet]
MachineMetadata=bastion

