# custom-units/blackhole.service
[Unit]
Description=block malicious ssh connections
After=sshd.service
Requires=sshd.service

[Service]
ExecStart=/bin/bash -c "while true ; do echo 'hello' ; sleep 1 done"
ExecStop=/bin/true

[X-Fleet]
MachineMetadata=bastion
