# custom-units/deis-backup.timer
[Unit]
Description=deis-backup-timer

[Timer]
OnCalendar=daily

[X-Fleet]
MachineOf=%p.service
