# custom-units/docker-cleanup.timer
[Unit]
Description=do a daily docker cleanup

[Timer]
OnCalendar=daily

[X-Fleet]
Global=true
