# custom-units/docker-volume-cleaner.timer
[Unit]
Description=clean the volumes daily

[Timer]
OnCalendar=daily

[X-Fleet]
Global=true
