# custom-units/image-cleaner.timer
[Unit]
Description=image-cleaner-timer

[Timer]
OnCalendar=daily

[X-Fleet]
Global=true
